# resources/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

su www-data
exec /usr/sbin/php-fpm7.4 --nodaemonize